# chef/templates/systemd.service.erb
[Unit]
Description=<%= @name %>

[Service]
Environment=<%= @environment.map {|key, val| %Q{"#{key}=#{val}"} }.join(' ') %>
ExecStart=<%= @command %>
KillSignal=<%= @stop_signal %>
User=<%= @user %>
WorkingDirectory=<%= @directory %>

[Install]
WantedBy=multi-user.target
